# roles/wayvnc/templates/wayvnc.service.j2
# wayVNC systemd user service
# Generated by Ansible - do not edit manually
# {{ ansible_facts['date_time']['iso8601'] }}

[Unit]
Description=wayVNC VNC server for Wayland
Documentation=https://github.com/any1/wayvnc
After=graphical-session.target
BindsTo=sway-session.target

[Service]
Type=simple
Environment=WAYLAND_DISPLAY={{ wayvnc_wayland_display }}
ExecStart=/usr/bin/wayvnc
Restart=on-failure
RestartSec=5

[Install]
WantedBy=sway-session.target
